# debian/postinst
#!/bin/sh
set -e

case "$1" in
    configure)
        # Create walpserver user and group if they don't exist
        if ! getent group walpserver >/dev/null; then
            addgroup --system walpserver
        fi
        if ! getent passwd walpserver >/dev/null; then
            adduser --system --ingroup walpserver --home /usr/share/walpserver walpserver
        fi

        # Create virtual environment
        python3 -m venv /usr/share/walpserver/venv

        # Activate virtual environment and install dependencies
        /usr/share/walpserver/venv/bin/pip install -r /usr/share/walpserver/requirements.txt

        # Create necessary directories
        mkdir -p /usr/share/walpserver/static/uploads
        mkdir -p /usr/share/walpserver/static/torrents
        mkdir -p /usr/share/walpserver/static/temp

        # Set permissions
        chown -R walpserver:walpserver /usr/share/walpserver
        chmod -R 755 /usr/share/walpserver

        # Enable and start service
        systemctl daemon-reload
        systemctl enable walpserver
        systemctl start walpserver
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

exit 0 